# opencode-deepseek4pro/vps-manager/systemd/vpsmgr-backup-prune.timer
[Unit]
Description=VPS Manager backup pruning
Documentation=https://github.com/vpsmgr/vpsmgr

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target