# opencode-bigpickle/vps-manager/systemd/vpsmgr-backup-prune.timer
[Unit]
Description=Daily VPS Manager backup pruning

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
