# opencode-kimi27/vps-manager/systemd/vpsmgr-backup-prune.timer
[Unit]
Description=Run VPS Manager backup pruning daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
